# roles/agama-client/templates/agama-client.service.j2
[Unit]
Description=Agama Client Service
After=network-online.target

[Service]
User=agama-client
ExecStart=/usr/local/bin/agama-client

[Install]
WantedBy=multi-user.target
